# dashboard/lightspeed-sync.service
[Unit]
Description=Lightspeed Customer Sync Service
After=network.target

[Service]
Type=simple
User=avi
WorkingDirectory=/home/avi/projects/mobile/rep/repair-dashboard
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/npx tsx /home/avi/projects/mobile/rep/repair-dashboard/scripts/sync-lightspeed-customers.ts
Restart=always
RestartSec=10
StandardOutput=append:/home/avi/projects/mobile/rep/repair-dashboard/logs/sync-out.log
StandardError=append:/home/avi/projects/mobile/rep/repair-dashboard/logs/sync-error.log

[Install]
WantedBy=multi-user.target
